# packaging/systemd/system-monitoring-dashboard.service
# Example systemd unit for running the monitoring dashboard outside Docker.
# Adjust paths as needed after building and installing the binary + web assets.

[Unit]
Description=System Monitoring Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Update WorkingDirectory if you install the web assets elsewhere.
WorkingDirectory=/opt/system-monitoring-dashboard
ExecStart=/usr/local/bin/dashboard
Environment=HOST_LABEL=
Environment=PORT=8080
Environment=WEB_ROOT=/opt/system-monitoring-dashboard/web
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target